# icinga2-ceb/content/etc/supervisor/conf.d/icinga2.conf
[supervisord]
nodaemon=true

[program:icinga2]
command=/bin/bash -c "/opt/supervisor/icinga2_supervisor"
redirect_stderr=true
